# meta-iotgw/recipes-iotgw/iotgwd/iotgwd/systemd/iotgwd.service
[Unit]
Description=IoT Gateway Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=notify

# Emplacements de config (override possibles)
Environment=IOTGWD_CONFIG=/etc/iotgw.yaml
Environment=IOTGWD_CONFDIR=/etc/iotgwd
EnvironmentFile=-/etc/default/iotgwd

# Lancement + reload
ExecStart=/usr/bin/iotgwd --config ${IOTGWD_CONFIG} --confdir ${IOTGWD_CONFDIR}
ExecReload=/bin/kill -HUP $MAINPID

# Watchdog (sd_notify WATCHDOG=1 depuis le démon)
WatchdogSec=30s

# Dirs gérés par systemd (créés au démarrage avec bons droits)
RuntimeDirectory=iotgwd
StateDirectory=iotgwd
# (optionnel)
#LogsDirectory=iotgwd
#CacheDirectory=iotgwd

# Politique de redémarrage
Restart=on-failure
RestartSec=2

# Durcissement raisonnable par défaut
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ProtectHome=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes

# Tell systemd this service has a config directory under /etc/iotgwd 
ConfigurationDirectory=iotgwd

# Écrire seulement là où nécessaire (StateDirectory ajoute déjà /var/lib/iotgwd
# aux chemins en écriture ; tu peux garder /etc/iotgwd si tu modifies des fragments)
ReadWritePaths=/etc/iotgwd

# Le service a besoin du réseau (sinon, remplace par une politique plus stricte)
IPAddressDeny=none

[Install]
WantedBy=multi-user.target
